# VPSApps/old/python/crontab
0 7 * * * /usr/bin/python3 /home/bs/Apps/lichessdigest.py
45 11 * * * /usr/bin/python3 /home/bs/Apps/pycomics.py
